# .config/supervisor/uwsgi.conf
[program:uwsgi]
command = /root/.pyenv/versions/deploy_eb_docker/bin/uwsgi --ini=config.uwsgi.uwsgi-app
